# EDCClient/cron/EDCClient.cron
# Update time from NTP ervery 3 hour
0 */3 * * * root /home/EDC/set_time.sh /home/EDC/edc_setup.conf &> /dev/null
# Compress back log every day
0 0 * * * root find /mnt/mmc/EDClog/ -type f -name "*.log" | xargs gzip &> /dev/null
# Delete backup log which is older than 7 Day every day
30 0 * * * root find /mnt/mmc/EDClog/ -type f -name "*.log.gz" -mtime +7 | xargs rm  &> /dev/null
* * * * * root /home/EDC/watchdog.sh  >> /home/EDCUpdate/EDCUpdate.log 2>&1
