# EDCClient/cron/EDCClient.cron
0 */3 * * * root /home/EDC/set_time.sh /home/EDC/edc_setup.conf &> /dev/null
